# etc/cron.d/00-rsync-gentoo-portage
 0  0  1  *   *		/opt/gentoo-rsync/rsync-gentoo-portage.sh
