# etc/cron.d/00-rsync-gentoo-portage
0 0 1 * * portage /opt/gentoo-rsync/rsync-gentoo-portage.sh
